# services/phosphor-chassis-power.service
[Unit]
Description=Chassis Power Monitor

[Service]
Restart=always
ExecStart=/usr/bin/phosphor-chassis-power

[Install]
WantedBy=multi-user.target
