# services/phosphor-power-control.service
[Unit]
Description=Phosphor Power Control

[Service]
Restart=always
ExecStart=/usr/bin/phosphor-power-control

[Install]
WantedBy=multi-user.target
