# services/phosphor-psu-monitor.service
[Unit]
Description=Phosphor Power Supply Monitor
Wants=obmc-mapper.target
After=obmc-mapper.target
Wants=mapper-wait@-org-openbmc-control-power0.service
After=mapper-wait@-org-openbmc-control-power0.service

[Service]
Restart=on-failure
ExecStart=phosphor-psu-monitor -c /usr/share/phosphor-psu-monitor/psu_config.json

[Install]
WantedBy=multi-user.target
